# provision/ansible/roles/ubuntu-setup/templates/borg.backup.timer.j2
[Unit]
Description=Borg User Backup Timer

[Timer]
WakeSystem=false
OnCalendar=daily

[Install]
WantedBy=timers.target
