# provision/ansible/roles/ubuntu-setup/templates/rclone.b2.service.j2
[Unit]
Description=Sync backup folder with B2

[Service]
Type=simple
ExecStart=/usr/bin/rclone sync {{ backup_pool }} secret:
User=root
Group=root

[Install]
WantedBy=timers.target
