# provision/ansible/roles/ubuntu-setup/templates/rclone.b2.timer.j2
[Unit]
Description=Perform an rclone sync daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
